# tools/creds-sync/com.todak.push-creds-to-tdcc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.todak.push-creds-to-tdcc</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/broneotodak/bin/push-to-tdcc.sh</string>
    <string>/tmp/tasp-creds-share.json</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/push-creds-to-tdcc.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/push-creds-to-tdcc.err</string>
</dict>
</plist>
