# tools/creds-sync/com.todak.slave-creds-push-tasp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.todak.slave-creds-push-tasp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/slave/bin/push-to-tasp.sh</string>
    <string>/tmp/slave-creds-share.json</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/slave-creds-push-tasp.out</string>

  <key>StandardErrorPath</key>
  <string>/tmp/slave-creds-push-tasp.err</string>
</dict>
</plist>
