# tools/nas-heartbeat/nas-heartbeat.service
[Unit]
Description=NAS Heartbeat to neo-brain
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/nas-heartbeat.sh
